# deploy/wardex.service
[Unit]
Description=Wardex XDR Security Agent
Documentation=https://github.com/pinkysworld/Wardex
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=wardex
Group=wardex
ExecStart=/usr/bin/wardex serve 8080 /usr/share/wardex/site
WorkingDirectory=/var/lib/wardex
Environment=WARDEX_CONFIG_PATH=/etc/wardex/wardex.toml
Environment=SENTINEL_CORS_ORIGIN=http://localhost:8080

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/wardex /var/log/wardex
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes
ProtectClock=yes
ProtectHostname=yes
RemoveIPC=yes
RestrictRealtime=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @mount
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE
PrivateDevices=yes
ProtectProc=invisible
ProcSubset=pid
UMask=0077
DevicePolicy=closed
IPAddressAllow=localhost
IPAddressAllow=10.0.0.0/8
IPAddressAllow=172.16.0.0/12
IPAddressAllow=192.168.0.0/16

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=2G
CPUQuota=200%

# Restart policy
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wardex

[Install]
WantedBy=multi-user.target
